# support/docker-smbload.cron
# alpine linux dcron version (goes to /var/spool/cron/crontabs/adup file)
# load from smb share
3 21 * * * /opt/adup/script/adup smbload
