# support/smbload.cron
# load from smb share
3 21 * * * adup /opt/adup/script/smbload
